# contrib/linux-packaging/archlinux/neucoind.service
[Unit]
Description=Nomiccoin daemon
After=network.target

[Service]
ExecStart=nomiccoind -logtimestamps

[Install]
WantedBy=default.target
